# tradingbot-monitor.service
[Unit]
Description=Trading Bot Monitor Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=tradingbot
Group=tradingbot
WorkingDirectory=/home/tradingbot/btc-trading
Environment=ENVIRONMENT=production
Environment=PYTHONPATH=/home/tradingbot/btc-trading
ExecStart=/usr/bin/python3 /home/tradingbot/btc-trading/scripts/bots/monitorBot.py --schedule 15m
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tradingbot /home/tradingbot/btc-trading/data

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
